# deploy/kse-nmt.service
# systemd service for KSE NMT Simulator
# Place at: /etc/systemd/system/kse-nmt.service
# Enable: sudo systemctl enable kse-nmt
# Start:  sudo systemctl start kse-nmt
# Logs:   sudo journalctl -u kse-nmt -f

[Unit]
Description=KSE NMT Exam Simulator
Documentation=https://github.com/kse/nmt-simulator
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/kse-nmt

# Environment file (copy .env.example → /opt/kse-nmt/.env and fill values)
EnvironmentFile=/opt/kse-nmt/.env

ExecStart=/usr/bin/node server.js

# Auto-restart on crash
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kse-nmt

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
